# qini/tun/etc/profile.d/path-sbin.sh
export PATH=/usr/local/sbin:/usr/sbin/:/sbin/:$PATH
