# qini/_config/autostart/clawsmail.desktop
[Desktop Entry]
Encoding=UTF-8
Type=Application
Name=claws-mail
Comment=
Exec=claws-mail
StartupNotify=false
Terminal=false

