# qini/_config/autostart/opera.desktop
[Desktop Entry]
Encoding=UTF-8
Type=Application
Name=opera
Comment=
Exec=opera
StartupNotify=false
Terminal=false

